# files/crond
#!/sbin/openrc-run

command="/usr/sbin/crond"
