# files/crond/crond
#!/sbin/openrc-run

description="crond - crontab runner"
supervisor="s6"
s6_service_path="${RC_SVCDIR}/s6-scan/${name}"

depend() {
   need s6-svscan
}
